# roles/wazuh/templates/wazuh-cert-renewer.service.j2
[Unit]
Description=Renew Wazuh certificates from Step CA
After=network.target

[Service]
Type=oneshot
ExecStart={{ wazuh_install_dir }}/scripts/renew-certs.sh

[Install]
WantedBy=multi-user.target
